# etc/profile.d/history.sh
#!/usr/bin/bash

HISTSIZE=500
HISTFILESIZE=0
HISTFILE=/dev/null
