# etc/systemd/system/checkupdates.timer
[Unit]
Description = Check Updates Status Timer

[Timer]
OnCalendar  = daily
AccuracySec = 15m

[Install]
WantedBy    = timers.target
